# WordPress/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command=/usr/sbin/php-fpm7 --nodaemonize

[program:nginx]
command= nginx -g "daemon off;"
